# deb-packaging/DEBIAN/postinst
#!/bin/sh
set -e

if ! id -g ophelia > /dev/null 2> /dev/null; then
    addgroup --system ophelia
fi

if ! id -u ophelia > /dev/null 2> /dev/null; then
    adduser --system --ingroup ophelia --home /var/lib/ophelia ophelia
fi

# chown -R ophelia:ophelia /var/lib/ophelia
# chmod 750 /var/lib/ophelia

systemctl enable ophelia-ci-server.service
systemctl start ophelia-ci-server.service
