# deb-packaging/DEBIAN/prerm
#!/bin/sh
systemctl stop ophelia-ci-server.service
systemctl disable ophelia-ci-server.service
rm -rf /var/lib/git/*
