# deb-packaging/DEBIAN/prerm
#!/bin/sh
set -e

systemctl stop ophelia-ci-server || true
systemctl disable ophelia-ci-server || true

if id -u ophelia &> /dev/null; then
    userdel -r ophelia
fi
